# lib/auto-pifi.service
[Unit]
Description=Automatic Pi Wifi Hotspot

[Service]
ExecStart=/usr/lib/systemd/scripts/auto-pifi
